# cron/crontab
# Main cron file. It runs all crons of all sites.
* * * * * root /scripts/servercrons.sh
